# tools/com.aptly.cloudflared.quick-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aptly.cloudflared.quick-tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/mac/bin/cloudflared-quick-tunnel.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/mac/Documents/Projects/Aptly</string>

  <key>StandardOutPath</key>
  <string>/tmp/aptly-cloudflared-launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/aptly-cloudflared-launchd.err.log</string>
</dict>
</plist>
